# dataserver/dataserver.service
[Unit]
Description=Wireless Sensor Mesh Dataserver (FastAPI + Catcher)
After=network.target nats-dataserver.service

[Service]
Type=simple
User=pmikol
WorkingDirectory=/home/pmikol/projects/wireless-sensor-mesh-daq/dataserver
Environment="PYTHONPATH=/home/pmikol/projects/wireless-sensor-mesh-daq/dataserver"
ExecStart=/usr/bin/env bash -c 'source .venv/bin/activate && exec python run_dataserver.py'

# Logging to journalctl
StandardOutput=journal
StandardError=journal

# Graceful shutdown behavior
KillSignal=SIGINT
KillMode=mixed
TimeoutStopSec=10

# Auto-restart on crash
Restart=on-failure
RestartSec=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
