# meshserver/nats-meshserver-internal.service
[Unit]
Description=NATS Server for MeshServer Internal PubSub (port 4222)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/nats-server -p 4222
User=nats
Group=nats

# Restart behavior
Restart=on-failure
RestartSec=2

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
CapabilityBoundingSet=
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
