# meshserver/nats-meshserver.service
[Unit]
Description=NATS Server for MeshServer
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/nats-server -p 4222
User=nats
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
